# ssh/rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start sshd service
# ==============================================================================

if ! bashio::var.has_value "$(bashio::app.port 22)"; then
    exec sleep infinity
fi

bashio::log.info "Starting the SSH daemon..."
exec /usr/sbin/sshd -D -e
